# rootfs/etc/profile.d/motd.sh
[ "$(id -u)" -eq 0 ] && echo 'HINT: try `abcde -o mp3:-V2 -x -p  [tracks]`'